# systemd/controller-mode.target
[Unit]
Description=Controller Mode Target

[Install]
WantedBy=multi-user.target